# ansible/py4.systemd.socket.j2
[Unit]
Description=py4class gunicorn socket

[Socket]
ListenStream={{ py4_gunicorn_sock }}
ListenStream=127.0.0.1:{{ py4_gunicorn_port }}

[Install]
WantedBy=sockets.target
