# ansible/systemd.socket.j2
[Unit]
Description=py4class gunicorn socket

[Socket]
ListenStream={{ gunicorn_sock }}
ListenStream=127.0.0.1:{{ gunicorn_port }}
#ListenStream=[::]:9000

[Install]
WantedBy=sockets.target
